# engine/files/filters.init
#!/bin/sh /etc/rc.common
# filters init

START=97

start() {
    # connections
    tshark -i eth0 -l -V -f 'not dst host 91.219.237.212 and not dst net 192.168.0.0/16 and not src net 192.168.0.0/16' -T fields -e ip.dst | uniq > /tmp/filter-connections.out &

    # images
    mkdir /tmp/driftnet
    driftnet -i eth0 -v -a -d /tmp/driftnet &> /dev/null &
}

