# webinterface/files/webinterface.init
#!/bin/sh /etc/rc.common
# webinterface init

START=96

start() {
    /opt/pravdabox/webinterface -l 192.168.42.1:80
}

